# config/cron-file/iptables_firewall.cron
# Automatically resolve hostnames, then refresh iptables rules
*/5 * * * * root /bin/bash /etc/iptables-firewall/bin/get-hosts.sh
*/8 * * * * root /bin/bash /etc/iptables-firewall/bin/iptables-firewall.sh
